# docker/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/with-contenv bash
exec /etc/gearcargo/scripts/postgres-run.sh
