# docker/rootfs/etc/s6-overlay/s6-rc.d/postgres-setup/up
/command/with-contenv /etc/gearcargo/scripts/postgres-setup.sh
