# docker/rootfs/etc/s6-overlay/s6-rc.d/migrate/up
/command/with-contenv /etc/gearcargo/scripts/migrate.sh
