# docker/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bash
exec /etc/gearcargo/scripts/redis-run.sh
