# docker/rootfs/etc/s6-overlay/s6-rc.d/init-perms/up
/command/with-contenv /etc/gearcargo/scripts/init-perms.sh
